# sparse/usr/lib/systemd/system/sensorfwd.service.d/s19mps.conf
[Service]
ExecStartPre=-/bin/sh /usr/bin/droid/sensorfw-pre-start.sh
